# docker/s6/s6-rc.d/patchdoll/run
#!/command/execlineb -P
with-contenv
cd /workspace
s6-setuidgid patchdoll
node /app/packages/core/dist/server.js
