# docker/s6/s6-rc.d/ngrok-tunnel/run
#!/command/execlineb -P
with-contenv
/bin/bash /etc/s6-overlay/s6-rc.d/ngrok-tunnel/run.sh
