# docker/s6/s6-rc.d/github-observe-runtime/up
with-contenv
/bin/bash /etc/s6-overlay/scripts/github-observe-runtime-check.sh
