# docker/s6/s6-rc.d/codex-auth/up
with-contenv
/bin/bash /etc/s6-overlay/s6-rc.d/codex-auth/auth.sh
