# docker/s6/s6-rc.d/slack-bridge/run
#!/command/execlineb -P
with-contenv
cd /app/packages/adapter-slack
s6-setuidgid patchdoll
node /app/packages/adapter-slack/app.cjs
